# packages/bsp/bpi-sunxi/usr/lib/systemd/system/bpi-hciuart.service
[Unit]
Description=Attach /dev/ttyAS1 to BlueZ stack
Wants=bluetooth.service
Before=bluetooth.service
After=dev-ttyAS1.device modprobe@.service

[Service]
Type=simple
ExecStart=/usr/local/bin/bpi-hciuart

[Install]
WantedBy=multi-user.target
